# alsa-volume-limiter.service
[Service]
ExecStart=/usr/bin/alsa-volume-limiter
ExecReload=/usr/bin/alsa-volume-limiter
Restart